# crontab
# 每天更新一次互联网时间
0  1  *  *  *  ntpdate -u time.nist.gov > /dev/null 2>&1
